# homebrew/Library/LaunchAgents/com.m-software-engineering.homebrew-maintenance.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "https://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.m-software-engineering.homebrew-maintenance</string>

	<key>ProgramArguments</key>
	<array>
		<string>/bin/zsh</string>
		<string>-lc</string>
		<string>$HOME/.config/homebrew/homebrew-maintenance.sh</string>
	</array>

	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>3</integer>
		<key>Minute</key>
		<integer>30</integer>
	</dict>

	<key>RunAtLoad</key>
	<false/>

	<key>LowPriorityIO</key>
	<true/>

	<key>Nice</key>
	<integer>10</integer>

	<key>StandardOutPath</key>
	<string>/tmp/com.m-software-engineering.homebrew-maintenance.out.log</string>

	<key>StandardErrorPath</key>
	<string>/tmp/com.m-software-engineering.homebrew-maintenance.err.log</string>
</dict>
</plist>
